# application/configuration/myproject.socket
[Unit]
Description=gunicorn socket for myproject

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target